# application-legacy.conf.in
description "Application Launching for Legacy Applications"
author "Ted Gould <ted@canonical.com>"

instance ${APP_ID}-${INSTANCE_ID}

start on application-legacy-start
stop on application-end or desktop-end

env APP_EXEC="echo Error"
env APP_EXEC_POLICY
env APP_URIS

pre-start script
	if pgrep -u system surfaceflinger ; then
		initctl set-env USING_SURFACE_FLINGER=1
	fi
end script

exec @pkglibexecdir@/desktop-exec "${APP_ID}" "${APP_URIS}"

post-start exec @pkglibexecdir@/zg-report-app open "application://${APP_ID}.desktop"
post-stop exec @pkglibexecdir@/zg-report-app close "application://${APP_ID}.desktop"
