# application-legacy.conf.in
description "Application Launching for Legacy Applications"
author "Ted Gould <ted@canonical.com>"

instance ${APP_ID}-${INSTANCE_ID}

start on application-legacy-start
stop on application-end or desktop-end

env APP_EXEC="echo Error"
env APP_EXEC_POLICY

pre-start exec @pkglibexecdir@/desktop-exec ${APP_ID}

script
	if [ -z $APP_EXEC_POLICY ]; then
		exec $APP_EXEC
	else
		exec aa-exec -p "$APP_EXEC_POLICY" -- $APP_EXEC
	fi
end script

post-start exec @pkglibexecdir@/zg-report-app open application://${APP_ID}.desktop
post-stop exec @pkglibexecdir@/zg-report-app close application://${APP_ID}.desktop
