# application-legacy.conf.in
description "Application Launching for Legacy Applications"
author "Ted Gould <ted@canonical.com>"

instance ${APP_ID}-${INSTANCE_ID}

start on application-legacy-start
stop on application-end or desktop-end

env APP_EXEC="echo Error"

pre-start script
	initctl set-env APP_EXEC="`@libexecdir@/desktop-exec ${APP_ID} \"${APP_URIS}\"`"
end script

exec $APP_EXEC

post-start exec @libexecdir@/zg-report-app open application://${APP_ID}.desktop
post-stop exec @libexecdir@/zg-report-app close application://${APP_ID}.desktop
