# application-click.conf.in
description "Application Launching for Click Applications"
author "Ted Gould <ted@canonical.com>"

instance ${APP_ID}

start on application-click-start
stop on application-end or desktop-end

env APP_ID
export APP_ID

env APP_URIS=""
env APP_DIR="/"

apparmor switch ${APP_ID}

env APP_EXEC="echo Error"

pre-start exec @libexecdir@/click-exec "${APP_ID}" "${APP_URIS}"

script
	cd ${APP_DIR}
	PATH=.:${PATH}
	exec ${APP_EXEC}
end script

post-start exec @pkglibexecdir@/zg-report-app open "application://${APP_ID}.desktop"
post-stop exec @pkglibexecdir@/zg-report-app close "application://${APP_ID}.desktop"
