# application-click.conf.in
description "Application Launching for Click Applications"
author "Ted Gould <ted@canonical.com>"

instance ${APP_ID}

start on application-click-start
stop on application-end or desktop-end

env APP_ID
env APP_EXEC="echo Error"
env APP_URIS
env APP_DIR
# For Surface Flinger
env APP_DESKTOP_FILE

env UPSTART_APP_LAUNCH_ARCH="@upstart_app_launch_arch@"
export UPSTART_APP_LAUNCH

apparmor switch ${APP_ID}

pre-start script
	if pgrep -u system surfaceflinger ; then
		initctl set-env USING_SURFACE_FLINGER=1
	fi

	exec @pkglibexecdir@/click-exec
end script

# Remember, this is confined
exec @pkglibexecdir@/exec-line-exec

# Disabled for Phone 1.0 which doesn't yet have ZG installed on it,
# so we don't need to log the data.
#
# post-start exec @pkglibexecdir@/zg-report-app open
# post-stop exec @pkglibexecdir@/zg-report-app close
