# upstart-jobs/application-failed.conf.in
description "Application Failing"

start on stopped application-legacy RESULT=failed or stopped application-click RESULT=failed or stopped application-snap RESULT=failed
task

exec ${SNAP}/@pkglibexecdir@/application-failed
