# upstart-jobs/application-legacy.conf.in
description "Application Launching for Legacy Applications"
author "Ted Gould <ted@canonical.com>"

instance ${APP_ID}-${INSTANCE_ID}

start on application-legacy-start
stop on application-end or desktop-end

env APP_ID
env APP_EXEC
env APP_EXEC_POLICY=""
env APP_URIS
# For Surface Flinger
env APP_DESKTOP_FILE

env APP_XMIR_ENABLE=1
export APP_XMIR_ENABLE

# This will be set to "unconfined" by desktop-exec if there is no confinement defined
apparmor switch $APP_EXEC_POLICY

pre-start script
	if pgrep -u system surfaceflinger ; then
		initctl set-env USING_SURFACE_FLINGER=1
	fi

	exec @pkglibexecdir@/desktop-exec
end script

exec @pkglibexecdir@/exec-line-exec

post-start exec @pkglibexecdir@/zg-report-app open
post-stop exec @pkglibexecdir@/zg-report-app close
