# upstart-jobs/application-click.conf.in
description "Application Launching for Click Applications"
author "Ted Gould <ted@canonical.com>"

instance ${APP_ID}

start on application-click-start
stop on application-end or desktop-end

env APP_ID
env APP_EXEC
env APP_URIS
env APP_DIR
# For Surface Flinger
env APP_DESKTOP_FILE
env APP_PROCESS_GROUP

env UPSTART_APP_LAUNCH_ARCH="@upstart_app_launch_arch@"
export UPSTART_APP_LAUNCH_ARCH

apparmor switch ${APP_ID}

pre-start script
	if pgrep -u system surfaceflinger ; then
		initctl set-env USING_SURFACE_FLINGER=1
	fi

	exec @pkglibexecdir@/click-exec
end script

# Remember, this is confined
exec @pkglibexecdir@/exec-line-exec

post-start script
	initctl setenv APP_PROCESS_GROUP=`initctl status | cut -d " " -f 4`
	@pkglibexecdir@/zg-report-app open
end script

post-stop script
	kill "-$(APP_PROCESS_GROUP)"
	@pkglibexecdir@/zg-report-app close
end script
