# upstart-jobs/application-legacy.conf.in
description "Application Launching for Legacy Applications"
author "Ted Gould <ted@canonical.com>"

instance ${APP_ID}-${INSTANCE_ID}

start on application-legacy-start
stop on application-end or desktop-end

env APP_ID
env APP_EXEC
env APP_EXEC_POLICY=""
env APP_URIS
env APP_DESKTOP_FILE_PATH
env APP_XMIR_ENABLE

# This will be set to "unconfined" by desktop-exec if there is no confinement defined
apparmor switch $APP_EXEC_POLICY
cgroup freezer

# Initial OOM Score
oom score 110

# This could be confined
exec ${SNAP}/@pkglibexecdir@/exec-line-exec

post-start exec ${SNAP}/@pkglibexecdir@/zg-report-app open
post-stop script
	${SNAP}/@pkglibexecdir@/zg-report-app close
	${SNAP}/@pkglibexecdir@/cgroup-reap-all

	DEVELOPER_MODE=`gdbus call --system --dest com.canonical.PropertyService --object-path /com/canonical/PropertyService --method com.canonical.PropertyService.GetProperty adb`
	if [ "$DEVELOPER_MODE" != "(true,)" ] ; then
		rm -f ${HOME}/.cache/upstart/application-legacy-${APP_ID}-${INSTANCE_ID}.log*
	fi
end script
