# upstart-jobs/untrusted-helper.conf.in
description "Untrusted Helpers installed from Click Packages"

start on untrusted-helper-start
stop on untrusted-helper-end or desktop-end

instance ${HELPER_TYPE}:${INSTANCE_ID}:${APP_ID}

env APP_ID
env APP_EXEC="echo Error"
env HELPER_TYPE
env INSTANCE_ID=""
env APP_URIS

env UBUNTU_APP_LAUNCH_ARCH="@ubuntu_app_launch_arch@"
env UBUNTU_APP_LAUNCH_LIBEXEC="${SNAP}/@pkglibexecdir@"

apparmor switch ${APP_ID}
cgroup freezer
oom score 800

# This is unconfined
pre-start script
	if [ -x "${UBUNTU_APP_LAUNCH_LIBEXEC}/${HELPER_TYPE}/exec-tool" ] ; then
		${UBUNTU_APP_LAUNCH_LIBEXEC}/${HELPER_TYPE}/exec-tool
	else
		echo "Unable to find exec tool for ${HELPER_TYPE}"
		exit -1
	fi
end script

# Remember, this is confined
exec ${UBUNTU_APP_LAUNCH_LIBEXEC}/exec-line-exec

post-stop exec ${UBUNTU_APP_LAUNCH_LIBEXEC}/cgroup-reap-all
